# scripts/com.newshub.build-news.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.newshub.build-news</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mituba/workdir/trash/20260131/scripts/cron-build-news.sh</string>
    </array>

    <!-- 6時間ごとに実行（スリープ復帰後にも対応） -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <!-- 標準出力/エラーのログ -->
    <key>StandardOutPath</key>
    <string>/Users/mituba/workdir/trash/20260131/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mituba/workdir/trash/20260131/logs/launchd-stderr.log</string>

    <!-- 環境変数 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- 起動時に実行しない -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
